# root/etc/s6-overlay/s6-rc.d/kavita/run
#!/usr/bin/with-contenv bash

if [ "$ATIVAK" = "1" ]; then
    cd /ativak
else
    cd /kavita
fi

cmd="./Kavita"
echo "*** executing => $cmd"
exec s6-setuidgid abc \
    $cmd
